# ansible/templates/supervisor.cron
@reboot	{{apache_user}}	{{supervisord_bin}} -c {{supervisord_conf}}
